# build/ansible/go-qcoordinator/files/go-qcoordinator.service
[Unit]
Description=Go Beanstalkd Queue Coordinator
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Restart=on-failure
ExecStart=/usr/local/bin/go-qcoordinator

[Install]
WantedBy=multi-user.target